# airootfs/etc/systemd/system/mados-ventoy-setup.service
[Unit]
Description=madOS Persistence Detection
DefaultDependencies=no
After=systemd-remount-fs.service local-fs.target
Before=mados-persistence-detect.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/mados-ventoy-setup.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
