# airootfs/etc/systemd/system/mados-persistence.service
[Unit]
Description=madOS Persistent Storage Setup
After=local-fs.target
Before=multi-user.target
ConditionPathExists=/run/archiso

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup-persistence.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
