# airootfs/etc/systemd/system/setup-claude-code.service
[Unit]
Description=Install Claude Code if not already present
After=network-online.target pacman-init.service
Wants=network-online.target
# Only run if claude binary is not found
ConditionPathExists=!/usr/bin/claude

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=HOME=/root
ExecStart=/usr/local/bin/setup-claude-code.sh
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
